# flagmodel-reranker.service
[Unit]
Description=FlagModel Reranker Service
Documentation=https://github.com/FlagAI-Open/FlagEmbedding
After=network.target

[Service]
Type=simple
User=li
Group=li

# 工作目录
WorkingDirectory=/data/FlagEmbedding

# ===== 服务配置 =====
# 只需修改此处 MODEL_NAME，MODEL_PATH 会自动生成
Environment="RERANKER_MODEL_NAME=BAAI/bge-reranker-v2-m3"
Environment="RERANKER_MODEL_PATH=/model/bge-reranker-v2-m3"
Environment="RERANKER_VENV_PATH=/data/FlagEmbedding/.venv"
Environment="PATH=/usr/bin:/usr/local/bin"

# 服务配置
Environment="RERANKER_HOST=0.0.0.0"
Environment="RERANKER_PORT=10020"
Environment="RERANKER_USE_FP16=false"
Environment="RERANKER_USE_GPU=false"

# Worker 数量
Environment="RERANKER_WORKERS=2"

# 启动命令
ExecStart=/data/FlagEmbedding/reranker-start.sh

# 重启策略
Restart=on-failure
RestartSec=10

# 资源限制
LimitNOFILE=65535
MemoryMax=32G
MemoryHigh=24G

# 安全性设置
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
